# edge/rootfs/etc/s6-overlay/s6-rc.d/init-seed/up
/usr/share/frameworks/dist/helmsman/helmsman seed-edge
